# docker/autopilot/tags/etc/supervisor.d/autopilot.ini
[program:autopilot]
process_name=%(program_name)s
command=/usr/bin/php /var/www/html/artisan horizon
autostart=true
autorestart=true
user=autopilot
redirect_stderr=true
stdout_logfile=/var/www/logs/autopilot.log
stopwaitsecs=3600