# docker/autopilot/tags/etc/supervisor.d/caddy.ini
[program:caddy]
process_name=%(program_name)s
command=caddy run --config "/etc/caddy/Caddyfile" --adapter "caddyfile"
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/logs/caddy.log
stopwaitsecs=3600