# docker/autopilot/tags/etc/supervisor.d/php.ini
[program:php-fpm]
process_name=%(program_name)s
command=/usr/sbin/php-fpm7 -F
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/logs/php-fpm.log
stopwaitsecs=3600