# docker/base/tags/etc/supervisor.d/alertmanager.ini
[program:alertmanager]
process_name=%(program_name)s
command=alertmanager --config.file=/etc/alertmanager/alertmanager.yml --web.external-url {AUTOPILOT_URL}/monitor/alertmanager/ --storage.path=/etc/alertmanager/data
autostart=true
autorestart=true
user=autopilot
redirect_stderr=true
stdout_logfile=/var/www/logs/alertmanager.log
stopwaitsecs=3600