# docker/base/tags/etc/supervisor.d/blackbox.ini
[program:blackbox_exporter]
process_name=%(program_name)s
command=blackbox_exporter --config.file=/etc/blackbox_exporter/blackbox_exporter.yml
autostart=true
autorestart=true
user=autopilot
redirect_stderr=true
stdout_logfile=/var/www/logs/blackbox_exporter.log
stopwaitsecs=3600