# docker/base/tags/etc/supervisor.d/cron.ini
[program:cron]
process_name=%(program_name)s
command=crond -f
autostart=true
autorestart=true
user=autopilot
redirect_stderr=true
stdout_logfile=/var/www/logs/cron.log
stopwaitsecs=3600